# systemd/mnt-terraControl.automount
[Unit]
Description=Automount for terraControl
ConditionPathExists=/mnt/terraControl

[Automount]
Where=/mnt/terraControl
TimeoutIdleSec=10

[Install]
WantedBy=multi-user.target
